# haos/lib/haos.service
[Unit]
Description=Home Assistant OS (KVM nested in LightOS)
Documentation=https://www.home-assistant.io/installation/linux
Wants=network-online.target
After=network-online.target
StartLimitIntervalSec=600
StartLimitBurst=5

[Service]
Type=exec
EnvironmentFile=/opt/haos/haos.conf
RuntimeDirectory=haos
RuntimeDirectoryMode=0750

ExecStartPre=/opt/haos/bin/haos-network.sh up
ExecStart=/opt/haos/bin/haos-launch.sh
ExecStop=/opt/haos/bin/haos-stop.sh
ExecStopPost=/opt/haos/bin/haos-network.sh down

KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=90s

Restart=on-failure
RestartSec=10s

ProtectSystem=strict
# /opt/haos holds runtime config, scripts, and qcow2/UEFI data (under data/).
# /run/haos is the RuntimeDirectory.
ReadWritePaths=/opt/haos /run/haos

[Install]
WantedBy=multi-user.target
